# 03.Django/02.newbie_study/sample.conf
[supervisord]
nodaemon=true
logfile=/app/sample.log

[program:counter]
command=python /app/counter.py
stdout_logfile=/app/counter.log
